# .config/sway/contrib/swayidle.initd
#!/sbin/openrc-run

command=swayidle
command_args="${SWAYIDLE_COMMAND}"

pidfile=/var/run/swayidle.pid
